# redhat-ddns/redhat-ddns.timer
[Unit]
Description=redhat-ddns timer

[Timer]
OnBootSec=5m
OnUnitInactiveSec=5m

[Install]
WantedBy=timers.target
